# osx/com.user.aria.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.aria</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/aria2c</string>
        <string>--enable-rpc</string>
        <string>-d /Users/mariano/Downloads</string>
    </array>
    <key>KeepAlive</key>
    <true/>
</dict>
</plist>